# docker/rtltcp.conf
[program:rtl_tcp]
command=rtl_tcp
autostart=true
autorestart=true
stderr_logfile=/var/log/rtl_tcp.err.log
stdout_logfile=/var/log/rtl_tcp.out.log
